# compose_config/pfss_manager/crontab
*/5 * * * * /app/sync_pfss.sh /var/lib/pfss /app
